# docker/leshan/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=java -jar /opt/leshan/leshan.jar

[program:client]
command=java -jar /opt/leshan/leshan-client-example.jar
