# docker/leshan/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=java -jar /opt/leshan/leshan.jar --prefix=/leshan

[program:client]
command=java -jar /opt/leshan/leshan-client-example.jar localhost 5683
